# deploy/chili-dashboard.service
[Unit]
Description=Chili Home Dashboard containers
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/chili-dashboard
# The installed Compose plugin targets API 1.52 while this Pi's Docker daemon
# supports API 1.41. Keep container startup compatible until Docker is upgraded.
Environment=DOCKER_API_VERSION=1.41
ExecStart=/usr/bin/docker compose -f compose.yaml -f compose.pi.yaml up -d --build
ExecStop=/usr/bin/docker compose -f compose.yaml -f compose.pi.yaml down
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
